# celery.ini
[program:celery]
command=celery -A bigdata_bi_test worker --loglevel=info
directory=/Users/hongjinquan/Documents/project/pythonProject/bigdata_bi_test
stdout_logfile=/Users/hongjinquan/Documents/project/pythonProject/bigdata_bi_test/logs/celery_out.log
stderr_logfile=/Users/hongjinquan/Documents/project/pythonProject/bigdata_bi_test/logs/celery_err.log
autorestart=true
redirect_stderr=true